# email-intelligence.service
[Unit]
Description=Email Intelligence Service - AI-Powered Email Management
Documentation=https://github.com/yourusername/email-intelligence
After=network.target

[Service]
Type=forking
User=yourusername
Group=yourusername
WorkingDirectory=/Users/iamomen/apple-mcp

# Environment variables
Environment="NODE_ENV=production"
Environment="BACKEND_PORT=8002"
Environment="FRONTEND_PORT=3000"

# Start command
ExecStart=/Users/iamomen/apple-mcp/email-service.sh start

# Stop command
ExecStop=/Users/iamomen/apple-mcp/email-service.sh stop

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:/Users/iamomen/apple-mcp/logs/service.log
StandardError=append:/Users/iamomen/apple-mcp/logs/service-error.log

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target